# server/service/auto-update-server.service
[Unit]
Description=Auto-update server code from Git on boot
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
User=danielarmengolaltayo
Group=danielarmengolaltayo
WorkingDirectory=/home/danielarmengolaltayo
ExecStart=/home/danielarmengolaltayo/Desktop/remote-toggle-module/server/scripts/auto_update_server.sh
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
