# server/service/internet-led.service
[Unit]
Description=Internet connectivity LED (server)
After=network-online.target
Wants=network-online.target
ConditionPathExists=/run/boot-ready

[Service]
Type=simple
User=root
Group=root
ExecStart=/usr/bin/python3 /home/pi/Desktop/remote-toggle-module/server/scripts/internet_led.py
Restart=on-failure
RestartSec=2
Environment=PYTHONUNBUFFERED=1

[Install]
WantedBy=multi-user.target
