# server/service/server-put-blink.path
[Unit]
Description=Trigger server-put-blink.service when state.json changes

[Path]
# Detecta reemplazos atómicos y escrituras
PathChanged=/home/pi/Desktop/remote-toggle-module/server/state.json
PathModified=/home/pi/Desktop/remote-toggle-module/server/state.json
Unit=server-put-blink.service

[Install]
WantedBy=multi-user.target
