# server/service/server-put-blink.path
[Unit]
Description=Trigger server-put-blink.service when state.json changes

[Path]
# Se dispara cuando cambia el archivo del estado del servidor
PathChanged=/home/pi/Desktop/remote-toggle-module/server/state.json

Unit=server-put-blink.service

[Install]
WantedBy=multi-user.target
