# client/service/auto-update-client.service
[Unit]
Description=Auto-update client code from Git on boot
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
User=pi
Group=pi
WorkingDirectory=/home/pi
ExecStart=/home/pi/Desktop/remote-toggle-module/client/scripts/auto_update_client.sh
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target