# legacy/both/service/internet-led.path
[Unit]
Description=Start internet-led when /run/boot-ready exists (client)

[Path]
PathExists=/run/boot-ready
Unit=internet-led.service

[Install]
WantedBy=multi-user.target
